# deploy/systemd/spotify-auth-backend.service
[Unit]
Description=Spotify Auth Backend (FastAPI)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=spotify-auth
Group=spotify-auth
WorkingDirectory=/opt/spotify-auth-backend/backend
Environment=PYTHONUNBUFFERED=1
Environment=HOST=127.0.0.1
Environment=PORT=8000
EnvironmentFile=/etc/spotify-auth-backend/backend.env
ExecStart=/opt/spotify-auth-backend/.venv/bin/uvicorn app.main:app --host ${HOST} --port ${PORT}
Restart=on-failure
RestartSec=5
TimeoutStartSec=30
TimeoutStopSec=30
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/var/log/spotify-auth-backend
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
